# overlays/firmware-debian/01-printer-data/root/usr/lib/systemd/system/init-printer-data.service
[Unit]
Description=Initialize printer data from /opt/lava/printer_data to /oem/printer_data
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/init-printer-data.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
